# printer-proxy.service
[Unit]
Description=Printer power-on proxy (port 9100)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/hubitat-printer-proxy/printer-proxy.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
